# sat-x.service
[Unit]
Description=sat-x System Monitoring Service
After=network.target

[Service]
# Assumes the user running this service is the same user who owns the project files
# User=your_username # Not needed for --user services

# Assumes a virtual environment named .venv in the project root
# Update this path if your venv or uvicorn location is different
ExecStart=/home/payload/sat-x/.venv/bin/uvicorn sat_x.main:app_instance --host 127.0.0.1 --port 8000

WorkingDirectory=/home/payload/sat-x
Restart=always
RestartSec=3
Environment="PYTHONPATH=/home/payload/sat-x/src"

# Standard output and error logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sat-x

[Install]
WantedBy=default.target
